# crontab/cron-jobs
*/1 * * * * bash -c "printenv && echo 'Hello World'" >> /var/log/cron.log 2>&1
