# crontab/cron-jobs
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
*/1 * * * * bash -c "printenv && echo 'Hello World'" >> /var/log/cron.log 2>&1
